# scripts/graceful-rm-cleanup.service
[Unit]
Description=Delete expired graceful-rm trash entries

[Service]
Type=oneshot
ExecStart=/usr/local/bin/graceful-rm --cleanup
